# docker/prod/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:run-start]
directory=/srv/frontend
command=yarn start
autostart=true
autorestart=true
